# caddy-proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/usr/bin/env bashio
bashio::log.info "Starting Caddy..."
exec /usr/bin/caddy run --config /data/config/caddy.json --adapter json 2>&1
