# caddy-proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Caddy Proxy Manager
# Service: Caddy
# ==============================================================================

# Fallback: create caddy.json if init-config didn't produce it
if ! bashio::fs.file_exists "/data/config/caddy.json"; then
    _http_port="${HTTP_PORT:-80}"
    bashio::log.info "caddy.json missing — creating default on port ${_http_port}"
    mkdir -p /data/config
    printf '{"apps":{"http":{"servers":{"srv0":{"listen":[":%s"],"routes":[]}}}}}\n' "${_http_port}" \
        > /data/config/caddy.json
fi

bashio::log.info "Starting Caddy..."
exec /usr/bin/caddy run --config /data/config/caddy.json 2>&1
