# caddy-proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/caddy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Caddy Proxy Manager
# Service: Caddy
# ==============================================================================
bashio::log.info "Starting Caddy..."
exec /usr/bin/caddy run --config /data/config/caddy.json --adapter json 2>&1
