# caddy-proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/init-data/run
#!/usr/bin/env bashio
bashio::log.info "Initializing /data directories..."
mkdir -p /data/db /data/certs /data/config /data/logs
bashio::log.info "Done."
