# caddy-proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/init-data/run
#!/usr/bin/env bashio
# shellcheck shell=bash
set +e
# ==============================================================================
# Home Assistant Add-on: Caddy Proxy Manager
# Initializing /data directories
# ==============================================================================
bashio::log.info "Initializing /data directories..."
mkdir -p /data/db /data/certs /data/config /data/logs
bashio::log.info "Done."