# caddy-proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/cpm/run
#!/usr/bin/with-contenv bashio

echo "--- CPM SERVICE STARTUP VERSION 1.4-ha.14 ---"
bashio::log.info "Starting caddy-proxy-manager..."

# Load environment variables manually to ensure Next.js sees them
# s6-overlay stores them in /run/s6/container_environment/ or /var/run/s6/container_environment/
S6_ENV_DIR="/var/run/s6/container_environment"
[ ! -d "$S6_ENV_DIR" ] && S6_ENV_DIR="/run/s6/container_environment"

if [ -d "$S6_ENV_DIR" ]; then
    for var in "$S6_ENV_DIR"/*; do
        export "$(basename "$var")=$(cat "$var")"
    done
fi

# Ensure SESSION_SECRET is also in .env file for Next.js as a fallback
# though process.env should be enough, some Next.js builds are picky.
# We write ALL environment variables to .env for maximum compatibility.
# In standalone mode, Next.js looks for .env in the root of the app.
{
    echo "SESSION_SECRET=$SESSION_SECRET"
    echo "ADMIN_USERNAME=$ADMIN_USERNAME"
    echo "ADMIN_PASSWORD=$ADMIN_PASSWORD"
    echo "DATABASE_URL=$DATABASE_URL"
    echo "CADDY_API_URL=$CADDY_API_URL"
    echo "PORT=$PORT"
    echo "BASE_URL=$BASE_URL"
    echo "NODE_ENV=production"
} > /app/.env

# Also create it in .env.local just in case
cp /app/.env /app/.env.local

bashio::log.info "Environment variables prepared. Starting node server.js..."

# Debug logs (FORCE INFO to be visible)
bashio::log.info "DEBUG: NODE_ENV=$NODE_ENV"
bashio::log.info "DEBUG: PORT=$PORT"
bashio::log.info "DEBUG: DATABASE_URL=$DATABASE_URL"
bashio::log.info "DEBUG: SESSION_SECRET length=${#SESSION_SECRET}"
if [ -n "$SESSION_SECRET" ]; then
    bashio::log.info "DEBUG: SESSION_SECRET starts with: ${SESSION_SECRET:0:4}..."
else
    bashio::log.err "DEBUG: SESSION_SECRET IS EMPTY!"
fi

cd /app
# Use env to explicitly pass all current exported variables to the node process
exec env SESSION_SECRET="$SESSION_SECRET" \
         ADMIN_USERNAME="$ADMIN_USERNAME" \
         ADMIN_PASSWORD="$ADMIN_PASSWORD" \
         DATABASE_URL="$DATABASE_URL" \
         CADDY_API_URL="$CADDY_API_URL" \
         PORT="$PORT" \
         BASE_URL="$BASE_URL" \
         NODE_ENV="production" \
         node server.js 2>&1
