# caddy-proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/cpm/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting caddy-proxy-manager..."

# Load environment variables manually to ensure Next.js sees them
# s6-overlay stores them in /var/run/s6/container_environment/
for var in /var/run/s6/container_environment/*; do
    export "$(basename "$var")=$(cat "$var")"
done

# Ensure SESSION_SECRET is also in .env file for Next.js as a fallback
# though process.env should be enough, some Next.js builds are picky.
echo "SESSION_SECRET=$SESSION_SECRET" > /app/.env

cd /app
exec node server.js 2>&1
