# caddy-proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/cpm/run
#!/usr/bin/env bashio
# shellcheck shell=bash
set +e
# ==============================================================================
# Home Assistant Add-on: Caddy Proxy Manager
# Service: Caddy Proxy Manager (Next.js)
# ==============================================================================
bashio::log.info "Starting caddy-proxy-manager..."

# Load environment variables from s6 container_environment
S6_ENV_DIR="/run/s6/container_environment"
[ ! -d "$S6_ENV_DIR" ] && S6_ENV_DIR="/var/run/s6/container_environment"
if [ -d "$S6_ENV_DIR" ]; then
    for var in "$S6_ENV_DIR"/*; do
        [ -f "$var" ] && export "$(basename "$var")=$(cat "$var")"
    done
fi

if [ -z "${SESSION_SECRET:-}" ]; then
    bashio::log.fatal "SESSION_SECRET is empty — init-config failed to write env vars"
    exit 1
fi

bashio::log.info "SESSION_SECRET length: ${#SESSION_SECRET}, PORT: ${PORT:-3000}"

{
    echo "SESSION_SECRET=${SESSION_SECRET}"
    echo "ADMIN_USERNAME=${ADMIN_USERNAME}"
    echo "ADMIN_PASSWORD=${ADMIN_PASSWORD}"
    echo "DATABASE_URL=${DATABASE_URL}"
    echo "CADDY_API_URL=${CADDY_API_URL}"
    echo "PORT=${PORT:-3000}"
    echo "NODE_ENV=production"
} > /app/.env
cp /app/.env /app/.env.local

cd /app
exec env SESSION_SECRET="${SESSION_SECRET}" \
         ADMIN_USERNAME="${ADMIN_USERNAME}" \
         ADMIN_PASSWORD="${ADMIN_PASSWORD}" \
         DATABASE_URL="${DATABASE_URL}" \
         CADDY_API_URL="${CADDY_API_URL}" \
         PORT="${PORT:-3000}" \
         NODE_ENV="production" \
         node server.js 2>&1
