# caddy-proxy-manager/rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Caddy..."
exec /usr/bin/caddy run --config /data/config/caddy.json --adapter json 2>&1
