# caddy-proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/init-config/run
#!/usr/bin/env bashio
# shellcheck shell=bash
set +e
# ==============================================================================
# Home Assistant Add-on: Caddy Proxy Manager
# Loading configuration and preparing environment
# ==============================================================================
bashio::log.info "--- ADD-ON CONFIGURATION ---"

SESSION_SECRET=$(bashio::config 'session_secret' 2>/dev/null || echo "")
if [ -z "${SESSION_SECRET}" ] || [ "${SESSION_SECRET}" = "null" ]; then
    if [ -f /data/.session_secret ]; then
        bashio::log.info "Loading session_secret from storage..."
        SESSION_SECRET=$(cat /data/.session_secret)
    else
        bashio::log.info "Generating random session_secret..."
        SESSION_SECRET=$(openssl rand -base64 32)
        echo "${SESSION_SECRET}" > /data/.session_secret
    fi
fi

ADMIN_USERNAME=$(bashio::config 'admin_username' 2>/dev/null || echo "admin")
ADMIN_PASSWORD=$(bashio::config 'admin_password' 2>/dev/null || echo "")
HTTP_PORT=$(bashio::config 'http_port' 2>/dev/null || echo "80")
HTTPS_PORT=$(bashio::config 'https_port' 2>/dev/null || echo "443")

[ -z "${ADMIN_USERNAME}" ] || [ "${ADMIN_USERNAME}" = "null" ] && ADMIN_USERNAME="admin"
[ -z "${HTTP_PORT}" ]      || [ "${HTTP_PORT}" = "null" ]      && HTTP_PORT="80"
[ -z "${HTTPS_PORT}" ]     || [ "${HTTPS_PORT}" = "null" ]     && HTTPS_PORT="443"

DATABASE_URL="file:/data/db/caddy-proxy-manager.db"
CADDY_API_URL="http://localhost:2019"
PORT="3000"
NODE_ENV="production"

# Write to s6 container_environment for longrun services
S6_ENV_DIR=""
for _c in /run/s6/container_environment /var/run/s6/container_environment; do
    [ -d "$_c" ] && { S6_ENV_DIR="$_c"; break; }
done

if [ -n "$S6_ENV_DIR" ]; then
    bashio::log.info "Writing env vars to ${S6_ENV_DIR}"
    printf "%s" "${SESSION_SECRET}"  > "$S6_ENV_DIR/SESSION_SECRET"
    printf "%s" "${ADMIN_USERNAME}"  > "$S6_ENV_DIR/ADMIN_USERNAME"
    printf "%s" "${ADMIN_PASSWORD}"  > "$S6_ENV_DIR/ADMIN_PASSWORD"
    printf "%s" "${HTTP_PORT}"       > "$S6_ENV_DIR/HTTP_PORT"
    printf "%s" "${HTTPS_PORT}"      > "$S6_ENV_DIR/HTTPS_PORT"
    printf "%s" "${DATABASE_URL}"    > "$S6_ENV_DIR/DATABASE_URL"
    printf "%s" "${CADDY_API_URL}"   > "$S6_ENV_DIR/CADDY_API_URL"
    printf "%s" "${PORT}"            > "$S6_ENV_DIR/PORT"
    printf "%s" "${NODE_ENV}"        > "$S6_ENV_DIR/NODE_ENV"
else
    bashio::log.warn "s6 container_environment not found (checked /run/s6 and /var/run/s6)"
fi

# Create default caddy.json if missing, using the configured HTTP port
if [ ! -f /data/config/caddy.json ]; then
    bashio::log.info "Creating default caddy.json on port ${HTTP_PORT}..."
    printf '{"apps":{"http":{"servers":{"srv0":{"listen":[":%s"],"routes":[]}}}}}\n' "${HTTP_PORT}" \
        > /data/config/caddy.json
fi

bashio::log.info "Configuration loaded: HTTP_PORT=${HTTP_PORT} HTTPS_PORT=${HTTPS_PORT} SECRET_LEN=${#SESSION_SECRET}"
