# caddy-proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/init-config/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Reading add-on options..."

SESSION_SECRET=$(bashio::config 'session_secret')
if bashio::var.is_empty "${SESSION_SECRET}"; then
    if [ ! -f /data/.session_secret ]; then
        bashio::log.info "Generating random session_secret..."
        SESSION_SECRET=$(openssl rand -base64 32)
        echo "${SESSION_SECRET}" > /data/.session_secret
    else
        bashio::log.info "Loading session_secret from storage..."
        SESSION_SECRET=$(cat /data/.session_secret)
    fi
fi

ADMIN_USERNAME=$(bashio::config 'admin_username')
ADMIN_PASSWORD=$(bashio::config 'admin_password')
HTTP_PORT=$(bashio::config 'http_port')
HTTPS_PORT=$(bashio::config 'https_port')

# Validate required fields
if [ ${#SESSION_SECRET} -lt 32 ]; then
  bashio::exit.nok "session_secret must be at least 32 characters long"
fi
if bashio::var.is_empty "${ADMIN_PASSWORD}"; then
  bashio::exit.nok "admin_password must not be empty"
fi

# Write env vars for s6 child processes
export SESSION_SECRET
export ADMIN_USERNAME
export ADMIN_PASSWORD
export HTTP_PORT
export HTTPS_PORT
export DATABASE_URL="file:/data/db/caddy-proxy-manager.db"
export CADDY_API_URL="http://localhost:2019"
export PORT="3000"

# Find s6 environment directory
S6_ENV_DIR="/var/run/s6/container_environment"
[ ! -d "$S6_ENV_DIR" ] && S6_ENV_DIR="/run/s6/container_environment"

if [ -d "$S6_ENV_DIR" ]; then
    printf "%s" "${SESSION_SECRET}"  > "$S6_ENV_DIR/SESSION_SECRET"
    printf "%s" "${ADMIN_USERNAME}"  > "$S6_ENV_DIR/ADMIN_USERNAME"
    printf "%s" "${ADMIN_PASSWORD}"  > "$S6_ENV_DIR/ADMIN_PASSWORD"
    printf "%s" "${HTTP_PORT}"       > "$S6_ENV_DIR/HTTP_PORT"
    printf "%s" "${HTTPS_PORT}"      > "$S6_ENV_DIR/HTTPS_PORT"
    printf "%s" "${DATABASE_URL}"    > "$S6_ENV_DIR/DATABASE_URL"
    printf "%s" "${CADDY_API_URL}"   > "$S6_ENV_DIR/CADDY_API_URL"
    printf "%s" "${PORT}"            > "$S6_ENV_DIR/PORT"
    printf "%s" "production"         > "$S6_ENV_DIR/NODE_ENV"
else
    bashio::log.warn "Could not find s6 environment directory!"
fi

# BASE_URL: use ingress URL if available, otherwise localhost
if bashio::var.has_value "$(bashio::addon.ingress_url 2>/dev/null)"; then
  BASE_URL=$(bashio::addon.ingress_url)
else
  BASE_URL="http://localhost:3000"
fi
export BASE_URL
[ -d "$S6_ENV_DIR" ] && printf "%s" "${BASE_URL}" > "$S6_ENV_DIR/BASE_URL"

# Next.js specific: sometimes it needs environment variables in a .env file or specifically formatted
# though server.js usually picks up process.env.
# We'll stick to s6 environment for now but ensure SESSION_SECRET is definitely there.

# Log values (masking secrets)
bashio::log.info "Configuration loaded."
bashio::log.info "SESSION_SECRET length: ${#SESSION_SECRET}"
bashio::log.info "ADMIN_USERNAME: ${ADMIN_USERNAME}"
bashio::log.info "PORT: ${PORT}"
bashio::log.info "DATABASE_URL: ${DATABASE_URL}"
